# ubuntu/services/Supervisor/my-ssr-app.conf
[program:views]
directory=/home/sava/apps/views/
command=npm run serve
user=sava
autostart=true
autorestart=true